# rice-vision.service
[Unit]
Description=Rice Vision Edge Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Change User to the user you SSH / log in as on the Pi
User=pi
WorkingDirectory=/home/pi/rice-vision/edge-client

# Runs the main orchestration script
ExecStart=/home/pi/rice-vision/edge-client/startup.sh

# Give Flask 60 s to come up before systemd marks the service failed
TimeoutStartSec=60

# Restart automatically if it crashes; wait 5 s between attempts
Restart=on-failure
RestartSec=5

# Pipe stdout/stderr to the system journal (journalctl -u rice-vision -f)
StandardOutput=journal
StandardError=journal

# Environment needed by Python venv
Environment="PATH=/home/pi/rice-vision/edge-client/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
